# deploy/dmbo-orchestrator.service
[Unit]
Description=DMBO Orchestrator
After=network-online.target redis-server.service
Wants=network-online.target
Requires=redis-server.service
StartLimitIntervalSec=60
StartLimitBurst=5

[Service]
Type=simple
EnvironmentFile=/etc/dmbo/dmbo.env
ExecStart=/usr/local/bin/dmbo-orchestrator
Restart=on-failure
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ProtectKernelTunables=true

[Install]
WantedBy=multi-user.target
